# salt/kube-scheduler/etc/supervisor/conf.d/kube-scheduler.conf
[program:kube-scheduler]
autorestart =  true
autostart   =  true
command     =  bash -c '. /etc/default/kube-scheduler;
                                  exec kube-scheduler ${ARGS}'
