# salt/kube-proxy/etc/supervisor/conf.d/kube-proxy.conf
[program:kube-proxy]
autorestart =  true
autostart   =  true
command     =  bash -c '. /etc/default/kube-proxy;
                                  exec kube-proxy ${ARGS}'
