# salt/mesos-master/etc/supervisor/conf.d/mesos-master.conf
[program:mesos-master]
autorestart =  true
autostart   =  true
command     =  mesos-init-wrapper master
