# salt/cadvisor/etc/supervisor/conf.d/cadvisor.conf
{% set minions = salt['roles.dict']('influxdb') -%}
[program:cadvisor]
autorestart =  true
autostart   =  true
{% if minions['influxdb'] -%}
command     =  cadvisor
                   -storage_driver=influxdb
                   -storage_driver_host={{ minions['influxdb'][0] }}:8086
                   -storage_driver_name=graphite
                   -storage_driver_user=root
                   -storage_driver_password=root
{% else -%}
command     =  cadvisor
{% endif -%}
