# salt/consul-server/etc/supervisor/conf.d/consul-server.conf
{% set nodename = salt['config.get']('nodename') -%}
[program:consul-server]
autorestart =  true
autostart   =  true
command     =  consul agent -data-dir /var/lib/consul
                            -node {{ nodename }}
                            -bootstrap-expect 1
                            -server
