# salt/mesos-dns/etc/supervisor/conf.d/mesos-dns.conf
[program:mesos-dns]
autorestart =  true
autostart   =  true
command     =  mesos-dns -config=/etc/mesos-dns/config.json
