# salt/mesos-slave/etc/supervisor/conf.d/mesos-slave.conf
[program:mesos-slave]
autorestart =  true
autostart   =  true
command     =  mesos-init-wrapper slave
