# salt/kubelet/etc/supervisor/conf.d/kubelet.conf
[program:kubelet]
autorestart =  true
autostart   =  true
command     =  bash -c '. /etc/default/kubelet;
                                  exec kubelet ${ARGS}'
