# salt/swift-object/etc/supervisor/conf.d/swift-object.conf
[program:swift-object]
autorestart =  true
autostart   =  true
command     =  swift-init object-server start
