# salt/fleet/etc/supervisor/conf.d/fleet.conf
[program:fleet]
autorestart =  true
autostart   =  true
command     =  fleetd
