# salt/cron/etc/supervisor/conf.d/cron.conf
[program:cron]
autorestart =  true
autostart   =  true
command     =  cron
