# salt/logio-harvester/etc/supervisor/conf.d/logio-harvester.conf
[program:logio-harvester]
autorestart =  true
autostart   =  true
command     =  log.io-harvester
