# salt/kubernetes-kubelet/etc/supervisor/conf.d/kubernetes-kubelet.conf
[program:kubernetes-kubelet]
autorestart =  true
autostart   =  true
command     =  bash -c '. /etc/default/kubernetes-kubelet;
                                  exec kubernetes-kubelet ${ARGS}'
