# salt/kubernetes-proxy/etc/supervisor/conf.d/kubernetes-proxy.conf
[program:kubernetes-proxy]
autorestart =  true
autostart   =  true
command     =  bash -c '. /etc/default/kubernetes-proxy;
                                  exec kubernetes-proxy ${ARGS}'
