# salt/flannel/etc/supervisor/conf.d/flannel.conf
[program:flannel]
autorestart =  true
autostart   =  true
command     =  flanneld
