# salt/aurora-scheduler/etc/supervisor/conf.d/aurora.conf
[program:aurora]
autorestart =  true
autostart   =  true
command     = /usr/local/aurora-scheduler/bin/schedule.sh
