# salt/lxd/etc/supervisor/conf.d/lxd.conf
[program:lxd]
autostart   =  true
autorestart =  true
command     =  lxd --tcp 0.0.0.0:8443
