# salt/kube-apiserver/etc/supervisor/conf.d/kube-apiserver.conf
[program:kube-apiserver]
autorestart =  true
autostart   =  true
command     =  bash -c '. /etc/default/kube-apiserver;
                                  exec kube-apiserver ${ARGS}'
