# salt/swift-object-expirer/etc/supervisor/conf.d/swift-object-expirer.conf
[program:swift-object-expirer]
autorestart =  true
autostart   =  true
command     =  swift-init object-expirer start
