# salt/logio-server/etc/supervisor/conf.d/logio-server.conf
[program:logio-server]
autorestart =  true
autostart   =  true
command     =  log.io-server
