# salt/etcd/etc/supervisor/conf.d/etcd.conf
{% set host = salt['config.get']('host') -%}
{% set ipv4 = salt['config.get']('fqdn_ip4') -%}
[program:etcd]
autorestart =  true
autostart   =  true
command     =  etcd \
                   -name      {{ host }} \
                   -addr      {{ ipv4[0] }}:4001 \
                   -peer-addr {{ ipv4[0] }}:7001
