# salt/consul-client/etc/supervisor/conf.d/consul-client.conf
{% set nodename = salt['config.get']('nodename') -%}
[program:consul-client]
autorestart =  true
autostart   =  true
command     =  consul agent -data-dir /var/lib/consul \
                            -node {{ nodename }}
