# salt/consul-web/etc/supervisor/conf.d/consul-web.conf
[program:consul-web]
autorestart =  true
autostart   =  true
command     =  consul agent -data-dir /var/lib/consul
                              -ui-dir /usr/share/consul/dist
