# salt/consul-web/etc/supervisor/conf.d/consul-web.conf
[program:consul-web]
autorestart =  true
autostart   =  true
command     =  consul agent -ui-dir /path/to/ui
