# salt/graphite-web/etc/supervisor/conf.d/graphite-web.conf
{% set py  = salt['config.get']('saltpath') + '/../graphite' -%}
{% set grp = salt['config.get']('graphite-web:group:name') -%}
{% set usr = salt['config.get']('graphite-web:user:name') -%}
[program:graphite-web]
autostart   =  true
autorestart =  true
command     =  gunicorn_django \
                   -u {{ usr }} \
                   -g {{ grp }} \
                   -b    unix:/var/run/gunicorn/graphite-web.sock \
                   --log-file=/var/log/gunicorn/graphite-web.log \
                      {{ py }}/settings.py
user        =  {{ usr }}
