# salt/graphite-web/etc/supervisor/conf.d/graphite-web.conf
{% from  'graphite-web/map.jinja'
   import graphite_web
   with   context -%}
{% set roles = [] -%}
{% set grp = graphite_web['group']['name'] -%}
{% set usr = graphite_web['user']['name'] -%}
{% set py  = salt['config.get']('saltpath') + '/../graphite' -%}
[program:graphite-web]
autostart   =  true
autorestart =  true
command     =  gunicorn_django \
                   -u {{ usr }} \
                   -g {{ grp }} \
                   -b    unix:/var/run/gunicorn/graphite-web.sock \
                   --log-file=/var/log/gunicorn/graphite-web.log \
                      {{ py }}/settings.py
user        =  {{ usr }}
